# proxmox.tjo.cloud/root/etc/systemd/system/folder2ram-sync.service
[Unit]
Description=folder2ram sync service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/folder2ram -syncall
