# ingress.tjo.cloud/root/systemd/system/dyndns.service
[Unit]
Description=Dynamic DNS Updater

[Service]
ExecStart=/usr/local/bin/dyndns

[Install]
WantedBy=multi-user.target
