# postgresql.tjo.cloud/ansible/root/etc/systemd/system/postgresql@.service.d/notify.conf
[Unit]
OnFailure=notify-webhook-failure@%n.service
