# shared/ansible/roles/lego/root/etc/systemd/system/lego-renew.timer
[Unit]
Description=Renew certificate using lego

[Timer]
OnCalendar=hourly
Persistent=true
Unit=lego-renew.service

[Install]
WantedBy=timers.target
