# v2.postgresql.tjo.cloud/ansible/root/etc/systemd/system/barman-backup.timer
[Unit]
Description=Backup PostgreSQL databases daily

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true
RandomizedDelaySec=4h

[Install]
WantedBy=timers.target
