# shared/ansible/roles/restic/root/etc/systemd/system/restic-backup.timer
[Unit]
Description=Restic Backup

[Timer]
OnCalendar=daily
Unit=restic-backup.service

[Install]
WantedBy=timers.target
