# proxmox.tjo.cloud/ansible/root/etc/systemd/system/folder2ram-sync.service
[Unit]
Description=folder2ram sync service

[Service]
Type=oneshot
ExecStart=/sbin/folder2ram -syncall
