# postgresql.tjo.cloud/root/etc/systemd/system/barman-cron.service
[Unit]
Description=Run Barman cron
Requires=postgresql.service
After=postgresql.service

[Service]
Type=simple
User=barman
ExecStart=/usr/bin/barman cron
