# mail.tjo.cloud/root/etc/systemd/system/restic-check.timer
[Unit]
Description=Restic Check

[Timer]
OnCalendar=daily
Unit=restic-check.service

[Install]
WantedBy=timers.target
